# etc/services.d/cron/run
#!/usr/bin/with-contenv bash

echo "starting cron service"
crontab -l
exec cron -f /ark/config/crontab
