# etc/services.d/cron/run
#!/usr/bin/with-contenv bash

echo "starting cron service"
exec cron -f -L 2 &
crondpid=$!

# load the crontab before we start
crontab /ark/config/crontab

wait $crondpid
