# etc/services.d/arkmanager/run
#!/usr/bin/with-contenv bash

exec arkmanager start --no-background --verbose
